# etc/xinetd.d/redirect-hudson-80
service redirect-hudson-80
{
        port                    = 80
        bind                    = 69.197.175.79
        type                    = UNLISTED
        disable                 = no
        socket_type             = stream
        wait                    = no
        user                    = nobody
        redirect                = 127.0.0.1 8080
}
